# memos/restic-memos-backup.timer
[Unit]
Description=Run memos backup daily at 22:00

[Timer]
OnCalendar=*-*-* 22:00:00
Persistent=true

[Install]
WantedBy=timers.target